# cron.d/custom
0 10 * * * php /opt/server_scripts/cron-payment-mail-reminder.php
